# ansible/roles/apt-discord/templates/apt-maintenance.timer.j2
[Unit]
Description=Run APT Maintenance Daily

[Timer]
OnCalendar=*-*-* {{ maintenance_time }}:00
Persistent=true

[Install]
WantedBy=timers.target
